# scripts/keepalive.cron
# Run every 2 minutes from any Linux host with cron installed.
# Replace python path and project path as needed.
*/2 * * * * cd /path/to/Ayurveda && /usr/bin/python3 scripts/keepalive_ping.py >> /tmp/aayu-keepalive.log 2>&1
